# Setup-Linux/init_d_arteranos.in
#!/bin/sh
# Start/stop the Arteranos dedicated server
#
### BEGIN INIT INFO
# Provides:          arteranos
# Required-Start:    $remote_fs $syslog $time
# Required-Stop:     $remote_fs $syslog $time
# Default-Start:     3 4 5
# Default-Stop:
# Short-Description: Arteranos dedicated server
### END INIT INFO

PATH=/bin:/usr/bin:/sbin:/usr/sbin
DESC="Arteranos dedicated server"
NAME=arteranos
DAEMON=/usr/bin/Arteranos-Server-Start
PIDFILE=/var/run/arteranos.pid
SCRIPTNAME=/etc/init.d/"$NAME"

test -f $DAEMON || exit 0

. /lib/lsb/init-functions


case "$1" in
start)	log_daemon_msg "Starting arteranos server" "arteranos"
        start-stop-daemon --start --quiet --oknodo --background --pidfile $PIDFILE --make-pidfile --exec $DAEMON 
        log_end_msg $?
	;;
stop)	log_daemon_msg "Stopping arteranos server" "arteranos"
        killall -u arteranos
        RETVAL=$?
        [ $RETVAL -eq 0 ] && [ -e "$PIDFILE" ] && rm -f $PIDFILE
        log_end_msg $RETVAL
        ;;
restart) log_daemon_msg "Restarting arteranos server" "arteranos" 
        $0 stop
        $0 start
        ;;
reload|force-reload) log_daemon_msg "Reloading configuration files for arteranos server" "arteranos"
        log_end_msg 0
        ;;
status)
        status_of_proc -p $PIDFILE $DAEMON $NAME && exit 0 || exit $?
        ;;
*)	log_action_msg "Usage: /etc/init.d/cron {start|stop|status|restart|reload|force-reload}"
        exit 2
        ;;
esac
exit 0
